# tomcat/tomcat8-jre7-nodejs7/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/supervisord_wrapper.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
